# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"
env DJANGO_SETTINGS_MODULE=tdd.settings.production

start on net-device-up
stop on shutdown

respawn

setuid ubuntu
chdir /home/ubuntu/sites/SITENAME/source/tdd
exec ../../virtualenv/bin/gunicorn --bind unix:/tmp/SITENAME.socket tdd.wsgi:application
